# os/linux/autostart/xmodmap.desktop
[Desktop Entry]
Name=pgupdn
Exec=xmodmap /home/rafa/.Xmodmap
Type=Application
