# linux/autostart/ma-launcher.desktop
[Desktop Entry]
Type=Application
Name=ma-launcher
Comment=create Manjaro-Architect launcher
Exec=/usr/bin/ma-launcher
